# packer/bootstrap/pyglidein/crontab
*/5     *   *   *   *   cd $HOME/pyglidein && ./submit-pyglideins.sh